# alpha/experiments/appview-infra/kit/tests/fixtures/audit_canonical_unbacked/generated/systemd/bad2-blob-0.service
[Service]
ExecStart=/usr/bin/rclone sync /var/lib/bad2/blobs/ r2:b/blobs
