# alpha/experiments/appview-infra/kit/config-templates/rclone.timer.tmpl
[Unit]
Description=${name} blob mirror timer (${blobdir})

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
