# alpha/experiments/appview-infra/kit/config-templates/service.service.tmpl
[Unit]
Description=${name} — appview-infra tenant service
Documentation=https://github.com/CroftCommunity/appview-infra
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${user}
Group=${user}
StateDirectory=${name}
StateDirectoryMode=0750
WorkingDirectory=${data_dir}
ExecStart=${exec_start}
Restart=always
RestartSec=2
# hardening (contract §4/§5)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ReadWritePaths=${data_dir}
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
