# alpha/experiments/appview-infra/kit/config-templates/snapshot.timer.tmpl
[Unit]
Description=${name} api snapshot timer (VACUUM INTO serve/)

[Timer]
OnBootSec=3min
OnUnitActiveSec=${snapshot_interval}
Persistent=true

[Install]
WantedBy=timers.target
